# automation/launchd/com.suneelworkspace.autolab.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.suneelworkspace.autolab</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/MAC/SuneelWorkSpace/autolab/scripts/autolab-watch</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>StandardOutPath</key>
  <string>/Users/MAC/SuneelWorkSpace/autolab/reports/launchd-autolab.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/MAC/SuneelWorkSpace/autolab/reports/launchd-autolab.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/MAC/SuneelWorkSpace</string>
</dict>
</plist>
